# jetson/systemd/bt-autoconnect.service
[Unit]
Description=Bluetooth auto connect (user)
After=graphical-session.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
ExecStart=%h/.local/bin/bt-autoconnect.sh
Restart=always
RestartSec=5
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
